# ansible/roles/actualizar_dns/files/actualizar_dns.timer
[Unit]
Description=Lanza actualizar_dns

[Timer]
OnCalendar=*-*-* *:30:00
Unit=actualizar_dns.service

[Install]
WantedBy=multi-user.target