# ansible/roles/owncloud/files/owncloud-data.mount
[Unit]
Description=Configuracion del montaje de la carpeta de datos de owncloud

[Mount]
Where=/usr/share/webapps/owncloud/data
What=minas.local.lan:/data0/Cloud
Type=nfs
DirectoryMode=0770
Options=rw,nfsvers=3
